# provision/system/sf.slice
[Unit]
Description=Slice for Segfault host containers
Before=slices.target

[Slice]
CPUAccounting=true
#CPUQuota=800%
#CPUWeight=50

IOAccounting=true
#IOWeight=100
