# guest/fs-root/etc/profile.d/l0pht.sh
# source'd during interactive shell login.

# Trampoline to this script:
[[ -e /usr/local/l0pht-guest/bin/l0pht-motd.sh ]] && /usr/local/l0pht-guest/bin/l0pht-motd.sh

